# controller-mk2/docker/service/jarvis-heartbeat-listener/run
#!/bin/bash
/sbin/setuser jarvis bash -c "cd ~ && exec bundle exec ./bin/jarvis-heartbeat-listener"
