# controller-mk2/docker/service/jarvis/run
#!/bin/bash
/sbin/setuser jarvis bash -c "cd ~ && exec bundle exec ./controller.rb"
